# lte_module/units/lte-ecm-bootstrap.service
[Unit]
Description=EC200U-CN ECM Mode Bootstrap
After=lte-enum-guard.service
Wants=lte-enum-guard.service
Before=systemd-networkd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /usr/local/lib/lte-module/ecm-bootstrap.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
